# php/protected/extensions/crontab/crontabs/my_crontab
*/70 * * * * php /home/tu/www/backend/php/yiicmd.php backup
